# package/skeleton-init-s6/skeleton/etc/s6-linux-init/build/run-image/service/s6-linux-init-shutdownd/run
#!/bin/execlineb -P

foreground { s6-rmrf fifo }
foreground { s6-mkfifo -m 0600 fifo }
s6-linux-init-shutdownd -c "/etc/s6-linux-init/current" -g 3000
