# package/skeleton-init-s6/skeleton/etc/s6-linux-init/build/run-image/service/s6-linux-init-shutdownd/run
#!/bin/execlineb -P
s6-linux-init-shutdownd -c "/etc/s6-linux-init/current" -g 3000
