# package/skeleton-init-s6/skeleton/etc/s6-linux-init/build/run-image/service/s6-linux-init-early-getty/run
#!/bin/execlineb -P
/sbin/getty 115200 ttyAMA0 vt100
