# package/skeleton-init-s6/skeleton/etc/s6-linux-init/build/run-image/service/s6-svscan-log/run
#!/bin/execlineb -P

fdmove -c 1 2
foreground { s6-rmrf fifo }
foreground { s6-mkfifo -m 0600 fifo }
redirfd -rnb 0 fifo
s6-log -bpd3 -- 1 t /run/uncaught-logs
