# package/skeleton-init-s6/skeleton/etc/s6-rc/service/mount-rwfs/up
mount -wt ext4 -o defaults,rw /dev/vda5 /mnt/rwfs
