# package/skeleton-init-s6/skeleton/etc/s6-rc/service/00/up
if { echo "init stage2 starting" }
hostname skabox
